# initramfs/etc/inittab
#ttyS0::respawn:/sbin/getty 38400 ttyS0
ttyS0::respawn:/bin/ash
